# scripts/mc-sync.timer
[Unit]
Description=Minecraft S3 Sync Timer

[Timer]
OnCalendar=*:0/5
Unit=mc-sync.service

[Install]
WantedBy=timers.target